# 70-paranoid.rules
SUBSYSTEM=="usb", ACTION=="add", RUN+="/usr/local/bin/paranoid-plug.sh"
SUBSYSTEM=="pci", ACTION=="add", RUN+="/usr/local/bin/paranoid-plug.sh"
